# tasks/backups/KausERPBackupDaily.timer
[Unit]
Description=Ejecutar KausERPBackupDaily automáticamente todos los días a las 19:00

[Timer]
OnCalendar=*-*-* 19:00:00
Persistent=true

[Install]
WantedBy=timers.target
